# consul-template.ini
[program:consul_template]
command=/usr/local/bin/consul-template -config /etc/consul-template.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
