# haproxy.ini
[program:haproxy]
command=/usr/local/bin/run-haproxy.sh
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
